# .config/systemd/user/offlineimap.service
[Unit]
Description=Offlineimap

[Service]
ExecStart=/usr/bin/offlineimap -u quiet
Restart=always

[Install]
WantedBy=default.target
