# .config/systemd/user/syncthing.service
[Unit]
Description=Syncthing

[Service]
ExecStart=/usr/bin/syncthing -no-browser
Restart=always

[Install]
WantedBy=default.target
